# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
user=root

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
priority=1                          ; ✅ Démarre en premier
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true                    ; ✅ Arrêt propre
killasgroup=true

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=2                          ; ✅ Démarre après PHP-FPM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true                    ; ✅ Arrêt propre
killasgroup=true
depends_on=php-fpm                  ; ✅ Attend PHP-FPM